# jsSeo-botfacing.conf
description "jsSeo bot facing server"

start on runlevel 2
stop on runlevel [!2]

chdir /var/www/jsSeo

#setuid web
#setgid web

respawn

exec /var/www/jsSeo/botfacing.py
